# contrib/systemd/ditox-prune.timer.template
[Unit]
Description=Run Ditox prune periodically

[Timer]
OnUnitActiveSec=@EVERY@
Persistent=true

[Install]
WantedBy=timers.target
