# contrib/systemd/ditox-sync.timer.template
[Unit]
Description=Run Ditox sync periodically

[Timer]
OnUnitActiveSec=@EVERY@
Persistent=true

[Install]
WantedBy=timers.target
