# contrib/systemd/ditox-prune.service
[Unit]
Description=Ditox prune

[Service]
Type=oneshot
ExecStart=%h/.cargo/bin/ditox-cli prune
WorkingDirectory=%h
Environment=RUST_LOG=warn
